# docker/quadlet/paperclip.container
[Unit]
Description=Paperclip AI Agent Orchestrator
Requires=paperclip-db.service
After=paperclip-db.service

[Container]
Image=paperclip-local
ContainerName=paperclip
Pod=paperclip.pod
Volume=%h/.local/share/paperclip:/paperclip:Z
Environment=HOST=0.0.0.0
Environment=BIZBOX_HOME=/paperclip
Environment=BIZBOX_DEPLOYMENT_MODE=authenticated
Environment=BIZBOX_DEPLOYMENT_EXPOSURE=private
Environment=BIZBOX_PUBLIC_URL=http://localhost:3100
EnvironmentFile=%h/.config/containers/systemd/paperclip.env

[Service]
Restart=on-failure
TimeoutStartSec=120

[Install]
WantedBy=default.target
